# deployment/ansible/roles/xplatform/files/rc.local
#!/bin/sh -e

su xplatform -c "/srv/run_platform"
